# crontab
* * * * * root openserp search
# New line character required!
